# service/burp-server/run
#!/bin/bash
exec /sbin/setuser abc /app/usr/sbin/burp -c /config/burp-server.conf -F
