# cron.d/bedup
MAILTO="$BURP_EMAIL"
# Put this file in /etc/cron.d and uncomment the following line to send a
# summary at 6 every morning.
# If your server is using a different config file to /etc/burp/burp-server.conf,
# change the first argument to its path.
#0 6 * * *       abc    /config/send_summary /config/burp-server.conf $BURP_EMAIL "Daily backup summary"

# The following will run file deduplication over all client storages every
# Saturday at 8 in the morning. Again, if your server is using a different
# config file to /etc/burp/burp-server.conf, change that argument.
0 8 * * 6	abc    /app/usr/sbin/bedup -l -c /config/burp-server.conf >> /app/burp-bedup 2>&1
